# etc/dnf-makecache.cron
#! /bin/bash

LOGFILE=/var/log/dnf-makecache.cron
STAMP=`date +%c`:
echo $STAMP >>$LOGFILE
/usr/bin/dnf makecache >>$LOGFILE
