# etc/dnf-makecache.cron
#! /bin/bash

/usr/bin/dnf makecache >>/var/log/dnf-makecache.cron
