# etc/yum.cron
#!/bin/sh

if [ -f /var/lock/subsys/yum ]; then
	/usr/sbin/yum -r 10 -e 0 -d 0 -y update yum
	/usr/sbin/yum -r 120 -e 0 -d 0 -y update
fi
