# etc/systemd/dnf-makecache.timer
[Unit]
Description=dnf makecache timer

[Timer]
OnBootSec=10min
OnUnitInactiveSec=1h
Unit=dnf-makecache.service

[Install]
WantedBy=basic.target
