# etc/systemd/dnf-automatic-download.timer
[Unit]
Description=dnf-automatic-download timer

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1d

[Install]
WantedBy=basic.target
